# home/frontear/.config/systemd/user/frontear-db.service
[Unit]
Description=Regenerate frontear package database

[Service]
Type=oneshot
ExecStart=%h/.local/bin/regen_db
